# root_alpine_pg15/etc/services.d/postgres/run
#!/usr/bin/with-contenv sh

echo "Starting postgres..."
s6-setuidgid postgres /usr/bin/postgres15 -D /config/postgres/
