# root_pg15/etc/services.d/postgres/run
#!/usr/bin/with-contenv sh

echo "Starting postgres..."
s6-setuidgid postgres /usr/lib/postgresql/15/bin/postgres -D /config/postgres/
